# dcmgr/contrib/upstart/vdc-nsa.conf

description "Wakame-VDC: naming service agent"
author "axsh Co."

start on runlevel [2345]
stop on runlevel [016]

respawn

normal exit 0

chdir /home/wakame/dcmgr/
exec ./bin/nsa -s amqp://localhost/ >> /var/log/nsa.log 2>&1
