# dcmgr/contrib/upstart/vdc-metdata.conf

description "Wakame-VDC: Dcmgr metadata server"
author "axsh Co."

start on runlevel [2345]
stop on runlevel [016]

respawn

normal exit 0

env RACK_ENV="production"

chdir /home/wakame/dcmgr/
#exec  bundle exec unicorn -c config/unicorn-metadata.conf web/metadata/config.ru >> /var/log/vdc-metadata.log 2>&1
exec  unicorn -c config/unicorn-metadata.conf web/metadata/config.ru >> /var/log/vdc-metadata.log 2>&1
