# dcmgr/contrib/upstart/vdc-collector.conf

description "Wakame-VDC: collector process"
author "axsh Co."

start on runlevel [2345]
stop on runlevel [016]

respawn

normal exit 0

chdir /home/wakame/dcmgr/
exec ./bin/collector -s amqp://localhost/ >> /var/log/collector.log 2>&1
