# debian/config/upstart/vdc-proxy.conf
description "Wakame-VDC: proxy server"
author "axsh Co."

env DAEMON=/usr/sbin/nginx

start on filesystem
respawn

pre-start script
    $DAEMON -t
    if [ $? -ne 0 ]
        then exit $?
    fi
end script
 
exec $DAEMON -g "daemon off;" -c /usr/share/axsh/wakame-vdc/tests/builder/conf/proxy.conf | tee /usr/share/axsh/wakame-vdc/tmp/vdc-proxy.log
