# hvc-hva/vendor/wakame/app_generators/wakame/templates/config/init.d/wakame-master
#!/bin/bash
# Start/stop the wakame master daemon.
#
### BEGIN INIT INFO
# Provides:          wakame-master
# Required-Start:    $local_fs $network rabbitmq-server
# Required-Stop:     $local_fs $network rabbitmq-server
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: 
# Description:       cron is a standard UNIX program that runs user-specified 
#                    programs at periodic scheduled times. vixie cron adds a 
#                    number of features to the basic UNIX cron, including better
#                    security and more powerful configuration options.
### END INIT INFO

USER=wakame
GROUP=wakame
PIDFILE=/var/run/wakame/wakame-master.pid

#LSBNAMES='-l'  # Uncomment for LSB name support in /etc/cron.d/

. /lib/lsb/init-functions
if [ -z $WAKAME_ROOT ]; then
    [ -f /etc/default/wakame ] && . /etc/default/wakame
    export GEM_HOME WAKAME_ROOT WAKAME_ENV WAKAME_CLUSTER_ENV
fi

start_master()
{
    curl -s -f --retry 3 http://169.254.169.254/2008-02-01/user-data | grep "node=agent" > /dev/null && {
	log_daemon_msg "Canceled to start master process for Wakame"
        exit 0
    }
    # Prepare the parent directory for the pidfile.
    typeset pid_dir=`dirname $PIDFILE`
    if su -c "test -w $pid_dir" $USER; then
      :
    else
      mkdir -p $pid_dir
      chown "$USER:$GROUP" $pid_dir
    fi
    
    start-stop-daemon --start --quiet --pidfile $PIDFILE --name wakame-master --startas $WAKAME_ROOT/bin/wakame-master -- -u $USER $LSBNAMES
    
}
stop_master()
{
    start-stop-daemon --stop --quiet --pidfile $PIDFILE --name wakame-master
}

case "$1" in
start)	log_daemon_msg "Starting master process for Wakame" "wakame-master"
	start_master
        log_end_msg $?
	;;
stop)	log_daemon_msg "Stopping master process for Wakame" "wakame-master"
	stop_master
        log_end_msg $?
        ;;
restart) log_daemon_msg "Restarting master process for Wakame" "wakame-master" 
	stop_master
	start_master
        log_end_msg $?
        ;;
reload|force-reload) log_daemon_msg "Reloading configuration files for " ""
        log_end_msg 0
        ;;
*)	log_action_msg "Usage: /etc/init.d/wakame-master {start|stop|restart|reload|force-reload}"
        exit 2
        ;;
esac
exit 0
