# debian/wakame-vdc.postinst
#!/bin/sh
# postinst script for wakame-vdc
#
# see: dh_installdeb(1)

set -e

prefix_path=/usr/share/axsh/wakame-vdc

if [ -x "/etc/init.d/wakame-vdc" ]; then
  update-rc.d wakame-vdc defaults >/dev/null
  if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
    invoke-rc.d wakame-vdc start || exit $?
  else
    /etc/init.d/wakame-vdc start || exit $?
  fi
fi

exit 0
