# contrib/etc/init/vdc-auth.conf
description "Wakame-VDC: auth server"
author "axsh Co."

start on started vdc-collector
stop on stopped vdc-collector

respawn
respawn limit 5 60

## rack name
env NAME=auth

script
    [ -f /etc/default/vdc-${NAME} ] && . /etc/default/vdc-${NAME}
    [ "x${RUN}" != "xyes" ] && {
      logger "[${NAME}] ${NAME} will not start. because RUN is not 'yes' in /etc/default/vdc-${NAME}."
      exit 0
    }

    ulimit -c ${DAEMON_COREFILE_LIMIT:-0}

    [ -d "${VDC_ROOT}" ] || {
      logger "no such directory: ${VDC_ROOT}"
      exit 1
    }
    cd ${VDC_ROOT}/frontend/dcmgr_gui

    exec bundle exec unicorn \
     -o ${BIND_ADDR:-127.0.0.1} \
     -p ${PORT:-3000} \
     -c ${UNICORN_CONF:-/etc/wakame-vdc/unicorn-common.conf} ./config.ru \
        >> /var/log/wakame-vdc/${NAME}.log 2>&1
end script
