# contrib/etc/init/vdc-dcmgr.conf
description "Wakame-VDC: api server"
author "axsh Co."

start on started vdc-collector
stop on stopped vdc-collector

respawn
respawn limit 5 60

## rack name
env NAME=dcmgr

script
    [ -f /etc/default/vdc-${NAME} ] && . /etc/default/vdc-${NAME}
    # Make RUN=yes effective only at auto start.
    [ -n "$UPSTART_EVENTS" -a "x${RUN}" != "xyes" ] && {
      logger "[${NAME}] Skip auto start for ${NAME}. Edit /etc/default/vdc-${NAME} to set RUN=yes."
      exit 0
    }

    ulimit -c ${DAEMON_COREFILE_LIMIT:-0}

    [ -d "${VDC_ROOT}" ] || {
      logger "no such directory: ${VDC_ROOT}"
      exit 1
    }
    cd ${VDC_ROOT}/dcmgr/

    exec bundle exec unicorn \
     -o ${BIND_ADDR:-0.0.0.0} \
     -p ${PORT:-9001} \
     -c ${UNICORN_CONF:-/etc/wakame-vdc/unicorn-common.conf} ./config-${NAME}.ru \
        >> /var/log/wakame-vdc/${NAME}.log 2>&1
end script
