# contrib/etc/init/vdc-fluentd.conf
description "Wakame-VDC: relay fluentd"
author "axsh Co."

start on vdc-net-device-up
stop on vdc-net-device-down
#start on (started rabbitmq-server)
#stop on (stopped rabbitmq-server or stopped mysql)

respawn
respawn limit 5 60

## agent name
env NAME=fluentd

script
    [ -f /etc/default/vdc-${NAME} ] && . /etc/default/vdc-${NAME}
    [ "x${RUN}" != "xyes" ] && {
      logger "[${NAME}] ${NAME} will not start. because RUN is not 'yes' in /etc/default/vdc-${NAME}."
      exit 0
    }

    ulimit -c ${DAEMON_COREFILE_LIMIT:-0}

    [ -d "${VDC_ROOT}" ] || {
      logger "no such directory: ${VDC_ROOT}"
      exit 1
    }
    cd ${VDC_ROOT}/vdc-fluentd/

    exec ${VDC_ROOT}/vdc-fluentd/bin/${NAME} \
        --config /etc/wakame-vdc/fluentd.conf \
        --log /var/log/wakame-vdc/${NAME}.log
end script
