# wakame-init/ubuntu/14.04/init/wakame-init-waiter.conf
# wakame-init-waiter - waiting for wakame-init.

description "waiting for wakame-init."

start on starting network-interface
instance $INTERFACE

task

env TIMEOUT=30

script
  wait_sec=${TIMEOUT}
  sleep_sec=3
  tries=0
  start_at=$(date +%s)
  while test -f /run/wakame-init.lock; do
    sleep ${sleep_sec}
    tries=$((${tries} + 1))
    if [ "$(($(date +%s) - ${start_at}))" -gt "${wait_sec}" ]; then
      echo "Retry Failure: Exceed ${wait_sec} sec: Retried ${tries} times" >&2
      break
    fi
  done

  touch /run/wakame-init.lock
  start wakame-init
  rm -f /run/wakame-init.lock
end script

