# upstart-scripts-used/misty/etc/init/vdc-hva.conf
description "Wakame-VDC: hypervisor agent"
author "axsh Co."

start on runlevel [2345]
stop on runlevel [016]
#start on started vdc-collector
#stop on stopped vdc-collector
start on started wakame-vdc
stop on stopped wakame-vdc

respawn
respawn limit 5 60

## agent name
env NAME=hva

script
    [ -f /etc/default/vdc-${NAME} ] && . /etc/default/vdc-${NAME}
    [ "x${RUN}" != "xyes" ] && {
      logger "[${NAME}] ${NAME} will not start. because RUN is not 'yes' in /etc/default/vdc-${NAME}."
      exit 0
    }

    ulimit -c ${DAEMON_COREFILE_LIMIT:-0}

    [ -d "${VDC_ROOT}" ] || {
      logger "no such directory: ${VDC_ROOT}"
      exit 1
    }
    cd ${VDC_ROOT}/dcmgr/

    amqp_server_uri=amqp://${AMQP_ADDR:-127.0.0.1}:${AMQP_PORT:-5672}/
    logger "[${NAME}] amqp_server_uri=${amqp_server_uri}"

    exec ./bin/${NAME} \
     -i ${NODE_ID:-`hostname`} \
     -s ${amqp_server_uri} 2>&1 \
      | /usr/bin/flog \
     -p /var/run/flog-vdc-${NAME}.pid \
        /var/log/vdc-${NAME}.log
end script
